# bunny-dns.service
[Unit]
Description=HestiaCP → Bunny DNS Sync Watcher
After=network-online.target hestia.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/hestia/plugins/bunny-dns/bunny-dns-watcher.sh
Restart=always
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bunny-dns
KillMode=process

[Install]
WantedBy=multi-user.target
